# meta-l4s-measure/recipes-measure/measure/files/measure.service
[Unit]
Description=Measuring service

[Service]
ExecStart=/usr/bin/measure

[Install]
WantedBy=multi-user.target