# recipes-core/scriptz/files/hello.service
[Unit]
Description=Hello service

[Service]
ExecStart=/usr/bin/hello-world

[Install]
WantedBy=multi-user.target